# deploy/feed-bot-webhook.service
[Unit]
Description=GitHub webhook listener for feed-bot auto-deploy
Documentation=https://github.com/adnanh/webhook
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy

# 仅监听 127.0.0.1，外部访问由 Nginx 反代 https://feed.fengshengxiu.club/deploy 进入
ExecStart=/usr/bin/webhook \
    -hooks /etc/feed-bot-webhook/hooks.json \
    -ip 127.0.0.1 \
    -port 9000 \
    -verbose

Restart=on-failure
RestartSec=5s

# 安全收紧
# 注意：不能加 NoNewPrivileges=yes —— 这条会让 sudo（setuid 二进制）
# 拒绝提权，导致 feed-bot-deploy.sh 里 `sudo systemctl restart pixiv-feed-bot`
# 失败并报「sudo: The "no new privileges" flag is set, which prevents sudo
# from running as root」。bot 自己的 service 不调 sudo，那里可以加。
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
# deploy 用户需要写仓库目录；其余路径无需写
ReadWritePaths=/opt/pixiv-feed-bot

# 日志走 journald；journalctl -u feed-bot-webhook 可看
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
